# powermonitor.service
[Unit]
Description=Power Monitor
After=network.target

[Service]
ExecStart=/home/pi/powermonitor.py
Type=simple

[Install]
WantedBy=default.target
